# docker/supervisord.conf
# =============================================================================
# Errora — supervisor config for the merged container.
#
# supervisord runs as PID-1's child (under tini) as root; each program drops to
# an unprivileged user. All logs go to the container's stdout/stderr so they are
# captured by `docker logs` / the platform's log driver.
#
# priority: lower starts first. App servers (10) before celery (20) before
# nginx (30) so the upstreams are listening before the proxy accepts traffic.
# Any program exiting unexpectedly takes the container down (exitcodes/startsecs)
# — let the orchestrator restart the whole unit rather than limp on.
# =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

# Shared log wiring (applied per program below; supervisor has no inheritance).
# stdout_logfile=/dev/fd/1 + maxbytes=0 forwards straight to the container.

# ---------------------------------------------------------------------------
# Django web — gunicorn (WSGI). Workers configurable via GUNICORN_WORKERS.
# ---------------------------------------------------------------------------
[program:gunicorn]
command=sh -c 'exec gunicorn errora.wsgi:application --bind 127.0.0.1:8000 --workers ${GUNICORN_WORKERS:-4} --timeout ${GUNICORN_TIMEOUT:-60} --access-logfile - --error-logfile -'
directory=/app/backend
user=errora
# supervisor switches uid but not HOME; without this gunicorn 26's control
# server tries /root/.gunicorn and hits EACCES. Point it at errora's home.
environment=HOME="/home/errora"
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=30
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

# ---------------------------------------------------------------------------
# Next.js standalone server.
# ---------------------------------------------------------------------------
[program:frontend]
command=node server.js
directory=/app/frontend
user=nextjs
# HOSTNAME must be 0.0.0.0: with 127.0.0.1 the standalone server's middleware
# self-proxy mishandles X-Forwarded-Proto and opens TLS to its own http port
# (EPROTO "wrong version number"). See Next.js standalone issue. nginx still
# reaches it via 127.0.0.1:3000 since 0.0.0.0 binds all interfaces.
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0",NEXT_TELEMETRY_DISABLED="1"
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=15
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

# ---------------------------------------------------------------------------
# Celery worker — all queues.
# ---------------------------------------------------------------------------
[program:celery-worker]
command=sh -c 'exec celery -A errora worker -Q ingest,ai,notifications,default -l info --concurrency=${CELERY_CONCURRENCY:-4}'
directory=/app/backend
user=errora
environment=HOME="/home/errora"
autostart=true
autorestart=true
startsecs=10
stopsignal=TERM
stopwaitsecs=60
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

# ---------------------------------------------------------------------------
# Celery beat — periodic scheduler. Schedule file lives in a writable tmp path.
# ---------------------------------------------------------------------------
[program:celery-beat]
command=celery -A errora beat -l info --schedule /tmp/celerybeat-schedule
directory=/app/backend
user=errora
environment=HOME="/home/errora"
autostart=true
autorestart=true
startsecs=10
stopsignal=TERM
stopwaitsecs=30
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

# ---------------------------------------------------------------------------
# Nginx — front proxy. Starts last; depends on the upstreams above.
# ---------------------------------------------------------------------------
[program:nginx]
command=nginx -g 'daemon off;'
user=root
autostart=true
autorestart=true
startsecs=5
stopsignal=QUIT
stopwaitsecs=15
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
